# scripts/event-monitor.service
[Unit]
Description=NGU Event Monitor Service
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/event-monitor
Environment=NODE_ENV=staging
Environment=AWS_REGION=us-east-1
Environment=DEBUG=*
Environment=NODE_OPTIONS="--trace-warnings --experimental-specifier-resolution=node"
EnvironmentFile=/home/ec2-user/event-monitor/.env
ExecStart=/usr/bin/node --experimental-specifier-resolution=node dist/run.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/event-monitor.log
StandardError=append:/var/log/event-monitor.error.log

# Ensure we have access to enough file descriptors
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target 